# debian-brski/brski-demo-app-deb/DEBIAN/postinst
#!/bin/bash

# Ask the user for a port number
echo "Please enter the port number for the demo-server service:"
read PORT

# Validate the port number
if ! [[ "$PORT" =~ ^[0-9]+$ ]] || [ "$PORT" -lt 1024 ] || [ "$PORT" -gt 65535 ]; then
    echo "Invalid port: $PORT. Please enter a number between 1024 and 65535."
    exit 1
fi

# Write the port number to the configuration file
echo "Port selected: $PORT"
mkdir -p /opt/demo-server/html  # Create the directory if it doesn't exist
echo "port=$PORT" > /opt/demo-server/html/server.conf

# Install and copy the brski tool and configs, respectively
wget -O /tmp/brski_0.2.3_arm64.deb https://github.com/nqminds/brski/releases/download/0.2.3/brski_0.2.3_arm64.deb
dpkg -i /tmp/brski_0.2.3_arm64.deb
cp /opt/demo-server/certs/* /etc/brski
rm /tmp/brski_0.2.3_arm64.deb

# Copy the systemd service file
cp /opt/demo-server/demo-server.service /etc/systemd/system/

# Reload systemd and enable the service
systemctl daemon-reload
systemctl enable demo-server.service
systemctl start demo-server.service
